# root/etc/systemd/system/shorewall.service.d/blacklist.conf
[Service]
ExecStartPre=/usr/share/nethserver-blacklist/load-ipsets
